# systemd/system/pncontrol.socket
[Unit]
Description=PixieNet XL Control Socket
Requires=pncontrol.service

[Socket]
ListenStream=0.0.0.0:9999

[Install]
WantedBy=sockets.target