# systemd/system/pncontrol.socket
[Unit]
Description=Pixie-Net XL Control Socket

[Socket]
ListenStream=0.0.0.0:9999
Accept=yes

[Install]
WantedBy=sockets.target
